# example.cron
## cronR job
## id:   abc
## tags: test1, test2
## desc: This is a really long description that will hopefully be wrapped
##   appropriately.
0 0 * * * testing!

# Minute   Hour   Day of Month       Month          Day of Week        Command    
# (0-59)  (0-23)     (1-31)    (1-12 or Jan-Dec)  (0-6 or Sun-Sat)                
0        2          12             *               0,6           /usr/bin/find
