# test.cron
## cronR job
## id:   abc
## tags: test1, test2
## desc: This is a really long description that will hopefully be wrapped
##   appropriately.
0 0 * * * testing!

## cronR job
## id:   123
## tags: test2, test3
## desc: 
0 0 * * * testing2!

 
